# devops-control.service
[Unit]
Description=DevOps Control Plane
Documentation=https://github.com/sauraku/devops-control
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple

# ── CHANGE THIS to match where you deploy the binary ──
WorkingDirectory=/opt/devops-control
ExecStart=/opt/devops-control/deploy/start.sh
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/opt/devops-control/data /opt/devops-control/.env.prod
PrivateTmp=yes

# Docker access requires these capabilities
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_FOWNER CAP_CHOWN CAP_SETGID CAP_SETUID CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_DAC_OVERRIDE CAP_FOWNER CAP_CHOWN CAP_SETGID CAP_SETUID CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
